# configurations/systemd/sol-watcher.path
[Path]
PathModified=/home/zero/solrunners

[Install]
WantedBy=multi-user.target