# configurations/systemd/sol.tate.service
[Unit]
Description=Springs of Life — Tate node (mpv X11/GPU video loop)
# Start after the storage.service has copied media to /storage.
After=storage.service network.target
Requires=storage.service

[Service]
Type=simple
User=zero

WorkingDirectory=/home/zero/solrunners

# Unbuffered Python output so journalctl shows print() in real time.
Environment=PYTHONUNBUFFERED=1

# xinit starts the bare X server (:0) and immediately runs xinit.tate.sh,
# which activates the venv and launches tate.py.
# -nocursor hides the mouse pointer during playback.
ExecStart=/usr/bin/xinit /home/zero/solrunners/configurations/xinit.tate.sh -- :0 -nocursor

Restart=on-failure
RestartSec=3

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
