# configurations/systemd/sol.entropy.ai.service
[Unit]
Description=Springs Of Life Entropy (AI, no WM)
After=network.target

[Service]
Type=simple
User=zero
WorkingDirectory=/home/zero/solrunners
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/xinit /home/zero/solrunners/configurations/xinit.entropy.ai.sh -- :0 -nocursor
Restart=always
RestartSec=3
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target

